# build/systemd/updatetracker.service
[Unit]
Description=Update tracker for MISP airgap components
After=network.target

[Service]
Type=simple
User=updatetracker
ExecStart=/usr/bin/python3 /opt/misp_airgap/build/update_tracker.py
Restart=on-failure

[Install]
WantedBy=multi-user.target